# initial_config/network-check.service
[Unit]
Description=Network Connectivity Check Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/network-check
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target